# dist/init/systemd/bashrpc.service
[Unit]
Description=bashrpc

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/bin/bashrpc -c /etc/bashrpc/bashrpc.yml

[Install]
WantedBy=multi-user.target
